# examples/supervisor-boot.conf
[program:boot]
user = root
command = /opt/cluster/etc/network.sh
autostart = true
autorestart = false
